# debian/wok.postrm
#!/bin/sh
set -e

if [ "$1" = "purge" ]; then
  update-rc.d wok remove >/dev/null
  rm -Rf /etc/wok
  rm -f /etc/nginx/sites-available/wok.conf
  invoke-rc.d nginx restart
fi

if [ "$1" = "remove" ]; then
  [ -L "/etc/nginx/sites-enabled/wok.conf" ] && rm -f /etc/nginx/sites-enabled/wok.conf
  invoke-rc.d nginx restart
fi

#DEBHELPER#
